# prototype/main/qt-style-config/qt-style-follow-lingmo/asset/overlay/etc/profile.d/qt-style.sh


#export QT_QPA_PLATFORMTHEME=lingmo
#export QT_PLATFORM_PLUGIN=lingmo
export QT_STYLE_OVERRIDE=lingmo

